# sssd.service.d/override.conf
[Service]
ExecStartPre=-/usr/bin/mkdir /var/log/sssd
